# packaging/napd.service
[Unit]
Description=napd — focus-aware power manager
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/lib/napd/napd.py
Restart=on-failure
RestartSec=3
# napd must never become the problem it watches for:
CPUWeight=10
Nice=10
MemoryHigh=64M

[Install]
WantedBy=default.target
