# packaging/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    systemctl --global disable napd.service >/dev/null 2>&1 || true
fi
exit 0
